# config/systemd/user/idle-lock.service
[[Unit]
Description=Idle Lock with Screenshot

[Service]
Type=oneshot
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/home/sagar/.python/scripts/idle_lock_snap.py

