# backup/home/sagar/.config/systemd/user/dotfile-backup.timer
[Unit]
Description=Run dotfile backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

