# backup/home/sagar/.config/systemd/user/dotfile-backup.service
[Unit]
Description=Backup dotfiles to GitHub

[Service]
Type=oneshot
ExecStart=%h/spdotfiles/backup_dotfiles.sh

